# crontab/cron.tab.templ
#run every 6 hours
0 */6 * * * BIN_DIR/load-all.sh >> CRONTAB_LOG_DIR/log.out
